# systemd/chia-poolboys.service.d/override.conf
[Service]
Environment="PATH=CHIA/venv/bin:POOL/venv/bin"