# systemd/chia-poolboys.service.d/override.conf
[Service]
Environment="PATH=/home/pieter/chia-blockchain/venv/bin"